# boards/RFSoC4x2/petalinux_bsp/meta-user/recipes-apps/rfclk-init/files/rfclk-init.service
[Unit]
Description=RFSoC4x2 LMK clock select/reset
Before=pynq-remote.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rfclk-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
